# crontab
*/5 * * * * /usr/local/bin/php /var/www/html/api/cron.php
